# systemd/epaper-dashboard.timer
[Unit]
Description=Run Dashboard Every 15 Minutes

[Timer]
OnBootSec=10sec
OnCalendar=*:*:0,15,30,45

[Install]
WantedBy=timers.target
